# etc/systemd/system/mysql_monitor_daemon.service
[Unit]
Description=MySQL Connection Monitor Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/mysql_monitor_daemon -s %i
User=root
Group=root

# Security hardening settings
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_SYS_RESOURCE
PrivateDevices=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallFilter=~@mount @swap @reboot @keyring @cpu-hotplug

# Reload the daemon when SIGHUP is received
ExecReload=/bin/kill -HUP $MAINPID

# Allow the daemon to write to the Prometheus metrics file
ReadWritePaths=/var/run/mysql_monitor.prom
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
